# probe-tech.service
[Unit]
Description=Probe Tech Control Web Interface
After=network.target

[Service]
Type=simple
User=pravar
WorkingDirectory=/home/pravar/probe-tech-control
ExecStart=/usr/bin/python3 spa_server.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
